# run_firefox.conf
description     "run firefox"
start on startup
task
exec /test-configs/open-firefox.sh
